# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:lita]
directory=/home/lita/sample-bot
command=lita

